# examples/git-webhook-receiver.service
[Unit]
Description=Git Webhook Receiver
After=network.target

[Service]
ExecStart=git-webhook-receiver.py

[Install]
WantedBy=multi-user.target
